# deploy/systemd/balancify-frontend.service
[Unit]
Description=HeiFam Frontend (Next.js)
After=network.target
Wants=network.target

[Service]
Type=simple
User=balancify
Group=balancify
WorkingDirectory=/opt/heifam/frontend
EnvironmentFile=/opt/heifam/frontend/.env.production
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3000
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target

